# aria2/root.h2o/etc/services.d/h2o/run
#!/usr/bin/execlineb -P

/usr/bin/h2o --mode master --conf /etc/h2o.conf
