# aria2/root/etc/services.d/lighttpd/run
#!/usr/bin/execlineb -S0
fdmove -c 2 1
/usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf
